# prometheus/supervisor.conf
[program:prometheus]
command=/usr/local/bin/prometheus --config.file=/etc/prometheus/prometheus.yml
killasgroup=true
stopasgroup=true
redirect_stderr=true
