# grafana/supervisor.conf
[program:grafana]
command= /usr/sbin/grafana-server --config /etc/grafana/grafana.ini --homepath /usr/share/grafana
killasgroup=true
stopasgroup=true
redirect_stderr=true
